# systemd/deadlock-server-picker.user.service
[Unit]
Description=Deadlock Server Picker - Apply firewall rules on login
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%h/.local/bin/deadlock-server-picker apply
ExecStop=%h/.local/bin/deadlock-server-picker reset

[Install]
WantedBy=default.target
